# system-templates/systemd/app.service.j2
# Systemd service file for {{ app_name }}
# Generated by WebOps

[Unit]
Description={{ app_name }} - WebOps Deployment
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User={{ webops_user }}
Group={{ webops_user }}
WorkingDirectory={{ repo_path }}

# Environment
Environment="PATH={{ venv_path }}/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH={{ repo_path }}"
{% for key, value in env_vars.items() %}
Environment="{{ key }}={{ value }}"
{% endfor %}

# Gunicorn configuration
ExecStart={{ venv_path }}/bin/gunicorn \
    --bind 127.0.0.1:{{ port }} \
    --workers {{ workers }} \
    --timeout 60 \
    --access-logfile {{ log_path }}/access.log \
    --error-logfile {{ log_path }}/error.log \
    --log-level info \
    {{ wsgi_module }}

# Restart policy
Restart=always
RestartSec=5

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ repo_path }} {{ log_path }}

[Install]
WantedBy=multi-user.target